# data_cycle/docker/log_generator.service
[Unit]
Description=Log Generator Service

[Service]
ExecStart=/usr/bin/python3 /path/to/log_generator.py
Restart=always

[Install]
WantedBy=multi-user.target
